# service/dogx.service
[Unit]
Description=OSdog X server

[Service]
Type=simple
ExecStart=/dog/dogx
EnvironmentFile=/etc/environment

[Install]
