# service/mb-af.service
[Unit]
Description=OSdog Motherboard Autofan Service

[Service]
Type=oneshot
RemainAfterExit=1
ExecStart=/dog/mb-af start
