# service/netmngr.service
[Unit]
Description=NetPacketManager

[Service]
Type=simple
EnvironmentFile=/dog/cfg/rig.cfg
ExecStart=/dog/sbin/netmngr $NETMANAGER_ARGS
