# service/sobakax.service
[Unit]
Description=Sobaka X server

[Service]
Type=simple
ExecStart=/sobaka/sobakax
EnvironmentFile=/etc/environment

[Install]
