# service/wd-kern.service
[Unit]
Description=OSdog Syslog Watchdog

[Service]
Type=simple
ExecStart=/dog/wd-kern
